# backend/almostprohibited-backend-indexer.service
[Unit]
Description=almostprohibited.ca backend indexer

[Service]
Type=oneshot
ExecStart={binary_path}/indexer

[Install]
WantedBy=multi-user.target
