# backend/almostprohibited-backend-indexer.timer
[Unit]
Description=Executes indexer daily

[Timer]
OnCalendar=*-*-* 0:00:00 America/Vancouver

[Install]
WantedBy=timers.target
